# deploy/ce-hub.service
[Unit]
Description=CE browser-node hub (WebSocket rendezvous + task dispatch + live stats)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ce-hub
ExecStart=/opt/ce-hub/ce-hub
Environment=CE_HUB_PORT=8970
Environment=CE_HUB_MODULES=/opt/ce-hub/modules
Environment=CE_HUB_DATA=/opt/ce-hub/data
# Storage limits so an anonymous caller can never fill this host's disk. The code defaults are
# already conservative (512 MiB / 200 apps / 100 domains / 500 db namespaces); these explicit
# values just make the relay's budget visible and bounded. Raise per host as trust grows.
Environment=CE_HUB_MAX_DATA_BYTES=1073741824
Environment=CE_HUB_MAX_APPS=500
Environment=CE_HUB_MAX_DOMAINS=200
Environment=CE_HUB_MAX_DB_NAMESPACES=1000
# Per-namespace IP rate limit (anti-spam) for the feedback app.
Environment=CE_HUB_RATELIMIT_NS=feedback
# Admin/moderation is IDENTITY-based, not a shared token: a request signed by one of these operator
# identities (owner id = sha256(pubkey)[..16]) is admin. An owner id is public, so it is safe here;
# the authority is the signature. This is the laptop's CE node identity (the operator).
Environment=CE_HUB_ADMIN_OWNER=52aa664acd1bde0a1da5224a037aea1a
Restart=always
RestartSec=2
# hardening: the hub holds no secrets and binds localhost only. It now persists app
# files, the KV database, and per-node uptime records under CE_HUB_DATA, so that path
# is explicitly made writable while the rest of the filesystem stays read-only.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ce-hub/data
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
